# scripts/memory-loop.service.example
# Memory Loop systemd user service
#
# Installation:
#   1. Copy to ~/.config/systemd/user/memory-loop.service
#   2. Edit paths and environment variables below
#   3. Run: systemctl --user daemon-reload
#   4. Run: systemctl --user enable --now memory-loop
#
# To start at boot without login:
#   loginctl enable-linger $USER
#
# Logs:
#   journalctl --user -u memory-loop -f

[Unit]
Description=Memory Loop - Obsidian vault interface
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USER/Projects/memory-loop

# Build frontend before starting (remove if you build separately)
ExecStartPre=/usr/bin/env bash -c 'cd frontend && bun run build'

# Start the backend server
ExecStart=/usr/bin/env bash -c 'cd backend && exec bun run start'

Restart=on-failure
RestartSec=5

# Environment - edit these for your setup
Environment=VAULTS_DIR=/path/to/your/vaults
Environment=PORT=3000
Environment=HOST=0.0.0.0

# Uncomment if bun is not in default PATH
#Environment=PATH=/home/YOUR_USER/.bun/bin:/usr/local/bin:/usr/bin:/bin

# Or use an environment file instead of inline variables:
#EnvironmentFile=/home/YOUR_USER/.config/memory-loop.env

[Install]
WantedBy=default.target
